# systemd/skydio-network-tester.service
[Unit]
Description=Skydio Network Readiness Tester
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/skydio-network-tester
Environment=PATH=/home/pi/skydio-network-tester/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/pi/skydio-network-tester/.venv/bin/python3 /home/pi/skydio-network-tester/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
